# system_files/usr/lib/systemd/user/copy-themes.service
[Unit]
Description=Copy system themes to ~/.themes

[Service]
Type=simple
ExecStart=/usr/libexec/copy-themes

[Install]
WantedBy=default.target