# system_files/usr/etc/profile.d/ydotool.sh
export YDOTOOL_SOCKET="/tmp/.ydotool_socket"